# builds/ubuntu0.0.9/sssh_server.service
[Unit]
Description=Service for sssh server

[Service]
Type=simple
ExecStart=/usr/local/bin/sssh_server